# D_WM-XFCE/install-scripts/systemd/archnexus-display-hotplug.path
[Unit]
Description=archnexus707 display hot-plug watcher

[Path]
PathChanged=/sys/class/drm
TriggerLimitBurst=20

[Install]
WantedBy=graphical-session.target
